# provision/supervisord.conf
[supervisord]
nodaemon=true

[program:node-inspector]
command=node-inspector --web-port=8081

[program:nodejs]
command=forever -w -c 'node --debug' app.js